# startup/blynk.sh
#!/bin/bash
#	Setting up Pi Jessie (removing pi-bluetooth and bluez to avoid install error
#	Updating and upgrading the system
#	Installing wiring PI..	
########################################################################
### BEGIN INIT INFO
# Provides:          autoupdater
# Required-Start:    $local_fs 
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: autoupdater
### END INIT INFO
sudo ./blynk --token=2f79d6db732040449192407b49934622
exit 1
